# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:geesefs]
command=/usr/local/bin/geesefs --no-checksum --memory-limit 4000 --max-flushers 32 --max-parallel-parts 32 --part-sizes 25 -f --endpoint %(ENV_AWS_ENDPOINT_URL_S3)s --list-type=2 --cache /data/geesefs-cache --disable-xattr --no-implicit-dir --no-dir-object %(ENV_BUCKET_NAME)s /mnt/data 
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pocker]
command=sh -c "ulimit -n 1000000 && /pocker"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
